# ansible/roles/prometheus_node_exporter/templates/tfc-meter-reading.timer
# {{ ansible_managed }}
[Unit]
Description=Timer to get current reading of our meters

[Timer]
OnCalendar=*:*:0/15

[Install]
WantedBy=timers.target
